# bootmachine/configuration/states/netcfg/netcfg.service.j2
.include /usr/lib/systemd/system/netcfg.service

[Unit]
After={{ pillar.netcfg_device }}.device
BindTo={{ pillar.netcfg_device }}.device
